# rootfs/etc/systemd/system/mydonglecloud-app.service
[Unit]
Description=MyDongle Cloud app service
After=systemd-modules-load.service
#After=network.target
#StartLimitIntervalSec=0

[Service]
Type=simple
Restart=oneshot
RestartSec=1
User=mdc
ExecStart=/home/mdc/app/app

[Install]
WantedBy=multi-user.target
