# rootfs/etc/systemd/system/mydonglecloud-init.service
[Unit]
Description=MyDongleCloud init
After=systemd-modules-load.service systemd-udevd.service
Requires=systemd-modules-load.service systemd-udevd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/modules/mydonglecloud/init.sh
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=sysinit.target
