# rootfs/etc/systemd/system/notify-failure@.service
[Unit]
Description=Failure Notification for %i

[Service]
Type=oneshot
ExecStart=/usr/local/modules/mydonglecloud/alert.sh "service" "%i"
