# rootfs/etc/update-motd.d/01-welcome
#!/bin/sh

cat /etc/issue.net
sh /usr/share/landscape/landscape-sysinfo.wrapper
